# board/opentrons/ot2/rootfs-overlay/etc/profile.d/ot-environ.sh
export RUNNING_ON_PI=1
export OT_SMOOTHIE_ID=AMA
