# board/opentrons/ot2/rootfs-overlay/etc/profile.d/ot-environ.sh
# Make sure that users of the api know this is a smoothie
export RUNNING_ON_PI=1
export OT_SMOOTHIE_ID=AMA
# Make sure scripts installed by pip are in the path
export PATH=$PATH:/var/user-packages/usr/bin
