# board/opentrons/ot2/rootfs-overlay/etc/systemd/system/syslog-ng.service.d/after-opentrons-update-server.conf
[Unit]
After=opentrons-update-server.service
