# board/opentrons/ot2/rootfs-overlay/etc/systemd/system/opentrons-gpio-setup.service
[Unit]
Description=Opentrons: Setup GPIO and blink
Before=systemd-growfs@var.service
DefaultDependencies=no
Conflicts=shutdown.target

[Service]
Type=simple
ExecStart=/usr/bin/ot-blinkenlights


[Install]
WantedBy=basic.target