# board/opentrons/ot2/rootfs-overlay/etc/systemd/system/opentrons.target
[Unit]
Description=Opentrons System
Requires=multi-user.target
Conflicts=emergency.target