# board/opentrons/ot2/rootfs-overlay/etc/systemd/system/opentrons-status-leds.service
[Unit]
Description=Opentrons: Run status LEDs
Before=opentrons-robot-server.service

[Service]
Type=simple
ExecStart=/usr/bin/ot-status-leds
Restart=on-failure
RestartSec=10s
RestartPreventExitStatus=10

[Install]
WantedBy=basic.target