# board/opentrons/ot2/rootfs-overlay/etc/systemd/system/opentrons-commit-machine-id.service
[Service]
Type=oneshot
ExecStart=/usr/bin/ot-commit-machine-id

[Install]
WantedBy=basic.target
