# board/opentrons/ot2/rootfs-overlay/lib/systemd/system/btusbstarter.service
[Unit]
Description=USB Bluetooth Adapter setup (on boot)
Requires=bluetooth.service
After=bluetooth.service
[Service]
Type=simple
ExecStart=/usr/bin/btusbstart
[Install]
WantedBy=opentrons.target
