# board/opentrons/ot2/rootfs-overlay/etc/systemd/system/jupyter-notebook.service.d/migrate-notebooks.conf
[Service]
ExecStartPre=/bin/sh /etc/systemd/system/jupyter-notebook.service.d/migrate-notebooks.sh