# board/opentrons/ot2/rootfs-overlay/etc/systemd/system/jupyter-notebook.service.d/migrate-notebooks.conf
[Service]
ExecStartPre=/usr/bin/env python /etc/systemd/system/jupyter-notebook.service.d/migrate-notebooks.py /var/resin-data user_storage/opentrons_data/jupyter $JUPYTER_NOTEBOOK_DIR -s