# board/opentrons/ot2/rootfs-overlay/etc/udev/rules.d/97-rtc-hctosys.rules
ACTION=="add", SUBSYSTEM=="rtc", RUN+="/sbin/hwclock --hctosys --utc"
